# nagios.conf
description "Nagios Upstart"

start on filesystem
stop on runlevel 5
respawn

# RUn Nagios
exec /usr/bin/nagios3 /etc/nagios3/nagios.cfg
